# app-utils/acme.sh/autobuild/overrides/usr/lib/systemd/user/acme-sh.timer
[Unit]
Description=acme.sh cron job timer

[Timer]
OnCalendar=daily
OnBootSec=30
Unit=acme-sh.service

[Install]
WantedBy=default.target
