# app-network/tftp-hpa/autobuild/overrides/usr/lib/systemd/system/tftpd.socket
[Unit]
Description=tftpd-hpa TFTP Server
Documentation=man:in.tftpd
Conflicts=tftpd.service

[Socket]
ListenDatagram=69

[Install]
WantedBy=sockets.target
