# app-network/wireless-regdb/autobuild/overrides/usr/lib/udev/rules.d/regulatory.rules
ACTION=="add" SUBSYSTEM=="module", DEVPATH=="/module/cfg80211", RUN+="/usr/bin/set-wireless-regdom"
